# scripts/systemd/stooq-downloader.service
[Unit]
Description=Stooq Data Downloader
After=network.target

[Service]
Type=oneshot
# Assuming the repository is cloned at ~/GitHub/stooq
# %h represents the user's home directory
WorkingDirectory=%h/GitHub/stooq
ExecStart=%h/GitHub/stooq/.venv/bin/python3 src/main.py
# Logs are sent to systemd journal (view with: journalctl --user -u stooq-downloader)
StandardOutput=journal
StandardError=journal
